# deploy/fluxcode.service
[Unit]
Description=FluxCode - AI API Gateway Platform
Documentation=https://github.com/DueGin/FluxCode
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=fluxcode
Group=fluxcode
WorkingDirectory=/opt/fluxcode
ExecStart=/opt/fluxcode/fluxcode
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fluxcode

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/fluxcode

# Environment - Server configuration
# Modify these values to change listen address and port
Environment=GIN_MODE=release
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080

[Install]
WantedBy=multi-user.target
